# root/etc/services.d/cron/run
#! /usr/bin/with-contenv bash

#启动trackers更新任务
exec crond -f -L /dev/null
