# root/etc/services.d/qBittorrent/run
#! /usr/bin/with-contenv bash

# 启动qBittorrent
exec qbittorrent-nox --webui-port=$WEBUIPORT   --profile=/config 
